# systemd/camofox.service
[Unit]
Description=Camofox Browser Server for Hermes Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=lincoln
WorkingDirectory=/home/lincoln/camofox-browser
Environment=CAMOFOX_PORT=9377
Environment=NODE_ENV=production
ExecStart=/home/lincoln/.nvm/versions/node/v24.14.0/bin/node server.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/lincoln/camofox-browser /tmp
PrivateTmp=true

[Install]
WantedBy=multi-user.target
