# systemd/hermes-agent.service
[Unit]
Description=Hermes Agent Telegram Gateway
After=network-online.target camofox.service
Wants=network-online.target

[Service]
Type=simple
User=lincoln
EnvironmentFile=/home/lincoln/.hermes/.env
ExecStart=/home/lincoln/.local/bin/hermes gateway start
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/lincoln /tmp
PrivateTmp=true

[Install]
WantedBy=multi-user.target
